# image_cleanup.service
[Unit]
Description=Image Cleanup Service

[Service]
Type=oneshot
ExecStart=/init/image_cleanup.sh
